# ssdpd.service.in
[Unit]
Description=Simple Service Discovery Protocol daemon (UPnP)
Documentation=man:ssdpd
Documentation=https://github.com/troglobit/ssdpd
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
EnvironmentFile=-@SYSCONFDIR@/default/ssdpd
ExecStart=/bin/ssdpd/ssdpd -n -w -i 900 -d http://%%s:1900/ssdpd.xml -T urn:txrx-com:device:DASCollector:1 -s eth0


[Install]
WantedBy=multi-user.target
